# etc/cron.d/shn-lc.conf
# download firmware
52   *    * * *   root    shipard-node fw-download
