# etc/cron.d/shn-mqtt.conf
# restart mqtt log
0 0    * * *   root      systemctl restart shn-mqtt-log
