# etc/systemd/shn-upload.service
[Unit]
Description=shipard node upload service ver 2.1
After=network.target auditd.service

[Service]
PIDFile=/var/run/shn-upload.pid
ExecStart=/usr/lib/shipard-node/tools/shn-upload.sh
Type=forking
Restart=on-failure

[Install]
WantedBy=multi-user.target
